# files/com.user.cfcalendarsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.cfcalendarsync</string>

    <!-- EDIT THE PATH BELOW to match where you placed cf_calendar_sync.py -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/bhavyadeephada/Desktop/cf x calender/files/cf_calender_sync.py</string>
    </array>

    <!-- Run once immediately when loaded, then every 24 hours -->
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>86400</integer>

    <key>StandardOutPath</key>
    <string>/tmp/cf_calendar_sync.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cf_calendar_sync.err.log</string>
</dict>
</plist>
